# build_files/kyth-duperemove.timer
[Unit]
Description=Weekly KythOS Steam/Wine dedupe pass

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=2h

[Install]
WantedBy=timers.target
