# build_files/kyth-topgrade-migrate.service
[Unit]
Description=KythOS migration: configure topgrade for KythOS (v3: system step, containers)
After=local-fs.target
ConditionPathExists=!/var/lib/kyth/migrations/topgrade-v3.done

[Service]
Type=oneshot
ExecStart=/usr/bin/kyth-topgrade-migrate

[Install]
WantedBy=multi-user.target
