# build_files/kyth-local-bin-migrate.service
[Unit]
Description=KythOS migration for old /usr/local helper paths
After=local-fs.target
ConditionPathExists=!/var/lib/kyth/migrations/local-bin-v1.done

[Service]
Type=oneshot
ExecStart=/usr/bin/kyth-local-bin-migrate

[Install]
WantedBy=multi-user.target
