# build_files/kyth-topgrade-migrate.service
[Unit]
Description=Kyth migration: configure topgrade for Kyth (disable distrobox, use sudo -n)
After=local-fs.target
ConditionPathExists=!/var/lib/kyth/migrations/topgrade-v2.done

[Service]
Type=oneshot
ExecStart=/usr/bin/kyth-topgrade-migrate

[Install]
WantedBy=multi-user.target
