# build_files/kyth-duperemove.service
[Unit]
Description=KythOS Steam/Wine dedupe pass
Documentation=https://github.com/mrtrick37/kyth
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/kyth-duperemove
Nice=19
IOSchedulingClass=idle
